# systemd/user/rechtsgebiet-sync.timer
[Unit]
Description=Nächtlicher Rechtsgebiets-Sync aus j-lawyer (04:10, nach dem 03:30-Fold)

[Timer]
OnCalendar=*-*-* 04:10:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
