# systemd/user/rm-backup.timer
[Unit]
Description=Nightly Rechtmaschine backup at 05:15

[Timer]
OnCalendar=*-*-* 05:15:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
