# systemd/user/rm-backup.service
[Unit]
Description=Rechtmaschine nightly backup: pg_dump + .env daily, volume tars Sundays, rsync to debian
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/var/opt/docker/rechtmaschine/scripts/rm_backup.sh
TimeoutStartSec=1h
Nice=15
